# deploy/com.yappa.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Yappa 游戏服务(launchd 常驻,登录即起,崩了自动拉起)。

  路径约定:文件里所有 /Users/YOUR_USERNAME 都是占位符。
  deploy/install.sh 在安装时会自动替换成你的家目录;
  如果你手动安装,请把每一处 /Users/YOUR_USERNAME 改成实际路径
  (仓库不在 ~/aiparty-exp 的话,WorkingDirectory 也要一起改)。

  绑定地址的取舍:这里【不加 −−lan】(XML 注释写不了连续两个半角横线,
  指的就是启动参数 lan 开关),即只绑 127.0.0.1。
  理由:Caddy 就在同一台机器上反代,走 127.0.0.1 完全够用;
  不绑 0.0.0.0 意味着 8747 端口不会暴露给局域网/公网,
  外界只能通过 Caddy 的 HTTPS 进来,这是安全上的正确姿势。
  入座链接的公网域名由 .env 里的 PUBLIC_BASE_URL 控制,与绑定地址无关。

  密钥:launchd 不读 shell 环境(.zshrc 的 export 无效),
  但程序自己会加载仓库根的 .env,所以这里【不需要】也【不应该】塞任何密钥。
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.yappa.server</string>

	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>-m</string>
		<string>modeb.simulator</string>
	</array>

	<!-- python -m 依赖 cwd 是仓库根;仓库不在这就改这里 -->
	<key>WorkingDirectory</key>
	<string>/Users/YOUR_USERNAME/aiparty-exp</string>

	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/Users/YOUR_USERNAME/Library/Logs/yappa/server.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/YOUR_USERNAME/Library/Logs/yappa/server.err.log</string>
</dict>
</plist>
